# unicorn/supervisord_unicorn.conf
[program:unicorn]
command=/root/exec_unicorn.sh
autostart=true
autorestart=true

